# deploy/supervisor.conf
[program:cccmu]
command=/opt/cccmu/bin/cccmu -p 8080 -l
directory=/opt/cccmu
user=root
startsecs=10
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/cccmu/cccmu.log